# waf/configs/coraza-forwarder.service
[Unit]
Description=Coraza Log Forwarder
After=network.target caddy.service

[Service]
Type=simple
User=caddy
Group=caddy
ExecStart=/usr/local/bin/coraza-forwarder \
    -coraza-log=/var/log/caddy/coraza_audit.log \
    -waf-log=/var/log/caddy/waf_wan.log \
    -api=http://localhost:8081/api \
    -poll=200ms
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/caddy
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target